# etc/systemd/system/main-ave.timer
[Unit]
Description=run main-ave.service

[Timer]
OnCalendar=13:10:00
Unit=main-ave.service

[Install]
WantedBy=timers.target