# etc/systemd/system/main-bel.service
[Unit]
Description=main bel

[Service]
Type=oneshot
ExecStart=/home/user1/main/main bel -src=v -srv=v -key=v -tag=v -mgn=v -mto=v
