# etc/systemd/system/main-foz.service
[Unit]
Description=main foz

[Service]
Type=oneshot
ExecStart=/home/user1/main/main foz -src=v -srv=v -key=v -tag=v -mgn=v -mto=v
