# etc/systemd/system/main-foz.service
[Unit]
Description=main foz

[Service]
Type=oneshot
ExecStart=/home/user1/main/main foz -pop='' -web='' -key='' -tag='' -mgn='' -mto=''
