# etc/systemd/system/main-foz.timer
[Unit]
Description=Run main-foz.service

[Timer]
OnActiveSec=2 hours
Unit=main-foz.service

[Install]
WantedBy=timers.target