# etc/systemd/system/main-ave.service
[Unit]
Description=main ave

[Service]
Type=oneshot
ExecStart=/home/user1/main/main ave -src=v -srv=v -key=v -tag=v -mgn=v -mto=v
