# etc/systemd/system/main-bel.service
[Unit]
Description=main bel

[Service]
Type=oneshot
ExecStart=/home/user1/main/main bel -ftp='' -web='' -key='' -tag='' -mgn='' -mto=''
