# etc/systemd/system/main-foz.timer
[Unit]
Description=run main-foz.service

[Timer]
OnCalendar=0*:15/2:00
Unit=main-foz.service

[Install]
WantedBy=timers.target