# etc/systemd/system/main-bel.service
[Unit]
Description=m15-worker bel

[Service]
Type=oneshot
ExecStart=/usr/bin/m15-worker bel -src=v -srv=v -key=v -tag=v -mgn=v -mto=v
User=m15
Group=m15
